# collector/crontab
# Каждый понедельник в 09:00 — снять статистику и сбросить
0 9 * * 1 cd /app && python collect.py >> /proc/1/fd/1 2>&1
